# service/gstreamer-qgc.service
[Unit]
Description=Pi Camera Video Stream to QGroundControl
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=robot
EnvironmentFile=/etc/eagle-stream.env
ExecStart=/home/robot/eagle-drone/scripts/stream.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
